# command-listener.service
[Unit]
Description=Daemon for listening for commands over the uart port
BindsTo=sys-devices-c280000.serial-tty-ttyTHS2.device
After=sys-devices-c280000.serial-tty-ttyTHS2.device
StartLimitIntervalSec=1
StartLimitBurst=3
FailureAction=poweroff
SuccessAction=poweroff

[Service]
Type=simple
User=root
WorkingDirectory=/home/nvidia/transfer-temp
ExecStart=/home/nvidia/payload-command-protocol/build/command-listener
Restart=on-failure

[Install]
WantedBy=multi-user.target
